# launchd/de.joerissen.journal-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.joerissen.journal-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/joerissen/ai/journal-bot/.venv/bin/python</string>
        <string>-m</string>
        <string>journal_bot.cli</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/joerissen/ai/journal-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- OPENROUTER_API_KEY kommt aus .env via python-dotenv,
             falls Du ihn dort hinterlegst. Alternativ hier eintragen. -->
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Jeden Montag 07:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/joerissen/ai/journal-bot/launchd/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joerissen/ai/journal-bot/launchd/stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
